# dockerfile/apps/nginx-rtmp/root/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh

# perform clean-up here

# use s6 to exit the container
# s6-svscanctl -t /var/run/s6/services